# otrs/etc/supervisord.d/otrs-slave.ini
[supervisord]
nodaemon=true
#loglevel=trace

[program:httpd]
command=/usr/sbin/httpd -DFOREGROUND
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true

[program:plugins]
command=/opt/otrs/bin/otrs.Console.pl Admin::Package::ReinstallAll 
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autostar=true
autorestart=false
user=otrs
